# systemd/pybackup_db.service
[Unit]
Description=Daily Mariadb database dump

[Service]
ExecStart=/usr/local/bin/pybackup.py -c db.json
